# systemd/hb-client.timer
[Unit]
Description=Heartbeat Monitor Client Timer

[Timer]
OnBootSec=30
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
