# systemd/hb-client.service
[Unit]
Description=Heartbeat Monitor Agent Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/heartbeat-monitor
Environment=CLIENT_CONFIG=/opt/heartbeat-monitor/config/client.yaml
Environment=PATH=/opt/heartbeat-monitor/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/heartbeat-monitor/.venv/bin/hb-agent daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
